# package/debian/postinst
#!/bin/sh
# postinst script for le-multiterminal

update-xorg-conf "Silicon.Motion" /etc/X11/xorg.conf.d/98-proinfo-*.conf

systemctl daemon-reload
systemctl enable le-multiterminal

# habilita e roda o xorg-daemon
systemctl enable xorg-daemon.socket
systemctl start xorg-daemon.socket

udevadm trigger

echo "Por favor reinicie o sistema para o multiterminal ser reconfigurado!"

exit 0
